# system_files/etc/systemd/system/mbp14-d3cold.service
[Unit]
Description=Disable NVMe d3cold for MBP14,3 suspend resume stability
After=local-fs.target
ConditionPathExists=/usr/local/libexec/mbp14-set-d3cold.sh

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/mbp14-set-d3cold.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
